# src/supervisord/conf/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:uwsgi]
command=uwsgi --ini /etc/uwsgi/mappening.ini --die-on-term
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
